# shoonya_service.service
[Unit]
Description=Shoonya Trading Platform - Execution Service
Documentation=file:///home/ec2-user/shoonya_platform/README.txt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/shoonya_platform

# Environment variables
Environment="PATH=/home/ec2-user/shoonya_platform/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"

# Service startup
ExecStart=/home/ec2-user/shoonya_platform/venv/bin/python /home/ec2-user/shoonya_platform/main.py

# ========================================
# ENHANCED RESTART POLICY
# ========================================
# Always restart on failure (including session recovery failures)
Restart=always
# Wait 10 seconds before restart (allows broker to stabilize)
RestartSec=10
# Allow up to 5 restarts in 5 minutes
StartLimitInterval=300
StartLimitBurst=5
# Exit code 1 (session failure) triggers restart
RestartForceExitStatus=1

# Resource limits
MemoryLimit=2G
CPUQuota=80%

# Process isolation
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/ec2-user/shoonya_platform/logs
ReadWritePaths=/home/ec2-user/shoonya_platform/shoonya_platform/persistence/data
ProtectHome=yes
ProtectKernelTunables=yes

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM
SendSIGKILL=yes

# Security hardening
OOMScoreAdjust=-100  # Prevent OOM killer
LimitNOFILE=65536
LimitNPROC=8192

# Logging to systemd journal (in addition to file logs)
StandardOutput=journal
StandardError=journal
StandardError=journal
SyslogIdentifier=shoonya-signal-processor

[Install]
WantedBy=multi-user.target
