# shoonya_service.service
[Unit]
Description=Shoonya Trading Platform - Unified Service
Documentation=https://github.com/yourusername/shoonya_platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/shoonya_platform

# ========================================
# ENVIRONMENT VARIABLES
# ========================================
Environment="PATH=/home/ec2-user/shoonya_platform/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"
Environment="PYTHONPATH=/home/ec2-user/shoonya_platform"

# ========================================
# SERVICE ENTRYPOINT
# ========================================
# Runs unified service (execution + dashboard + option supervisor)
ExecStart=/home/ec2-user/shoonya_platform/venv/bin/python /home/ec2-user/shoonya_platform/main.py

# ========================================
# RESTART POLICY
# ========================================
# Always restart on failure (including session recovery failures)
Restart=always
# Wait 10 seconds before restart (allows broker to stabilize)
RestartSec=10
# Allow up to 5 restarts in 5 minutes
StartLimitInterval=300
StartLimitBurst=5
# Exit code 1 (session failure) triggers restart
RestartForceExitStatus=1

# ========================================
# RESOURCE LIMITS
# ========================================
# Memory limit (prevents runaway memory usage)
MemoryLimit=2G
# CPU limit (80% of one core)
CPUQuota=80%

# ========================================
# SECURITY & ISOLATION
# ========================================
# Private /tmp directory
PrivateTmp=yes
# Prevent privilege escalation
NoNewPrivileges=yes
# Protect system files
ProtectSystem=strict
# Allow writes only to logs and data directories
ReadWritePaths=/home/ec2-user/shoonya_platform/logs
ReadWritePaths=/home/ec2-user/shoonya_platform/shoonya_platform/persistence/data
ReadWritePaths=/home/ec2-user/shoonya_platform/shoonya_platform/market_data/option_chain/data
# Protect home directories
ProtectHome=yes
# Protect kernel parameters
ProtectKernelTunables=yes

# ========================================
# GRACEFUL SHUTDOWN
# ========================================
# 30-second timeout for graceful shutdown (matches main.py shutdown logic)
TimeoutStopSec=30
# Send SIGTERM first (triggers signal_handler in main.py)
KillSignal=SIGTERM
# Force kill if SIGTERM doesn't work
SendSIGKILL=yes

# ========================================
# PROCESS LIMITS
# ========================================
# Prevent OOM killer from targeting this service
OOMScoreAdjust=-100
# File descriptor limit
LimitNOFILE=65536
# Process limit
LimitNPROC=8192

# ========================================
# LOGGING
# ========================================
# Log to systemd journal (in addition to file logs in logs/ directory)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shoonya-platform

[Install]
WantedBy=multi-user.target

# ========================================
# INSTALLATION INSTRUCTIONS
# ========================================
# 1. Copy this file to systemd:
#    sudo cp ~/shoonya_platform/shoonya_service.service \
#        /etc/systemd/system/shoonya_platform.service
#
# 2. Reload systemd:
#    sudo systemctl daemon-reload
#
# 3. Enable auto-start on boot:
#    sudo systemctl enable shoonya_platform
#
# 4. Start service:
#    sudo systemctl start shoonya_platform
#
# 5. Check status:
#    sudo systemctl status shoonya_platform
#
# For detailed instructions, see: SERVICE_INSTALLATION_LINUX.md
# ========================================
