# deployment/systemd/shoonya_stop.timer
[Unit]
Description=Shoonya Platform Auto-Stop Timer (Daily 12:00 AM)
Documentation=Stop trading service after market hours

[Timer]
# Stop at midnight daily
OnCalendar=*-*-* 00:00:00
# Allow manual stop
Persistent=true
# Accuracy of 1 minute is fine for shutdown
AccuracySec=1min

[Install]
WantedBy=timers.target
