# deployment/dnss.service
[Unit]
Description=DNSS (Delta Neutral Short Strangle) Strategy Service
After=network.target

[Service]
Type=simple
User=trading
WorkingDirectory=/opt/shoonya_platform
Environment="PYTHONUNBUFFERED=1"
Environment="DASHBOARD_ENV=primary"
Environment="DNSS_CONFIG=/opt/shoonya_platform/shoonya_platform/strategy_runner/saved_configs/dnss_nifty_weekly.json"
EnvironmentFile=/opt/shoonya_platform/config_env/primary.env

ExecStart=/opt/shoonya_platform/venv/bin/python main.py

# Restart on failure
Restart=on-failure
RestartSec=10

# Graceful shutdown
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dnss

[Install]
WantedBy=multi-user.target

