# deployment/trading.service
[Unit]
Description=Trading Service
Documentation=https://github.com/yourusername/trading
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/shoonya

# ========================================
# ENVIRONMENT VARIABLES
# ========================================
Environment="PATH=/opt/shoonya/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"
Environment="PYTHONPATH=/opt/shoonya"

# ========================================
# SERVICE ENTRYPOINT
# ========================================
ExecStart=/opt/shoonya/venv/bin/python /opt/shoonya/main.py

# ========================================
# RESTART POLICY
# ========================================
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5
RestartForceExitStatus=1

# ========================================
# RESOURCE LIMITS
# ========================================
MemoryLimit=2G
CPUQuota=80%

# ========================================
# SECURITY & ISOLATION
# ========================================
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/shoonya/logs
ReadWritePaths=/opt/shoonya/shoonya_platform/persistence/data
ReadWritePaths=/opt/shoonya/shoonya_platform/market_data/option_chain/data
ProtectHome=yes
ProtectKernelTunables=yes

# ========================================
# GRACEFUL SHUTDOWN
# ========================================
TimeoutStopSec=30
KillSignal=SIGTERM
SendSIGKILL=yes

# ========================================
# PROCESS LIMITS
# ========================================
OOMScoreAdjust=-100
LimitNOFILE=65536
LimitNPROC=8192

# ========================================
# LOGGING
# ========================================
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading

[Install]
WantedBy=multi-user.target
