# utilities/deployment/systemd/trading_stop.timer
[Unit]
Description=Trading Platform Auto-Stop Timer (Mon-Fri midnight IST)
Documentation=https://github.com/GauravKomarewar/TRADING_TERMINAL

[Timer]
# Stop at midnight Monday-Friday (covers all markets: NSE/NFO 15:30, MCX 23:30/23:55)
OnCalendar=Tue-Sat *-*-* 00:00:00
Timezone=Asia/Kolkata
Unit=trading_stop.service
# Catch up and fire if the system was down during scheduled time
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
