# utilities/deployment/master.service
[Unit]
Description=Shoonya Master Account Manager
Documentation=https://github.com/GauravKomarewar/TRADING_TERMINAL
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/shoonya_platform

Environment="PATH=/home/ubuntu/shoonya_platform/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu/shoonya_platform"

ExecStart=/home/ubuntu/shoonya_platform/venv/bin/python /home/ubuntu/shoonya_platform/master_manager.py --env config_env/master.env

Restart=always
RestartSec=5
MemoryMax=512M
CPUQuota=20%

NoNewPrivileges=yes
ProtectSystem=no
ProtectHome=no

TimeoutStopSec=15
KillSignal=SIGTERM
SendSIGKILL=yes
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-master

[Install]
WantedBy=multi-user.target
