# utilities/deployment/trading.service
[Unit]
Description=Trading Service (primary)
Documentation=https://github.com/GauravKomarewar/TRADING_TERMINAL
After=network-online.target
Wants=network-online.target
After=gateway.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/shoonya_platform

Environment="PATH=/home/ubuntu/shoonya_platform/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"
Environment="PYTHONPATH=/home/ubuntu/shoonya_platform"

ExecStart=/home/ubuntu/shoonya_platform/venv/bin/python /home/ubuntu/shoonya_platform/main.py --env config_env/primary.env

Restart=always
RestartSec=10
RestartForceExitStatus=1

MemoryMax=2G
CPUQuota=80%

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=no
ProtectHome=no
ProtectKernelTunables=yes

TimeoutStopSec=30
KillSignal=SIGTERM
SendSIGKILL=yes

OOMScoreAdjust=-500
LimitNOFILE=65536
LimitNPROC=8192

StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading

[Install]
WantedBy=multi-user.target
